# helloworld.service
[Unit]
Description="Hello World Script"

[Service]
ExecStart=PATH/helloworld_script.sh
