# helloworld.service
[Unit]
Description="Hello World Script"

[Service]
ExecStart=$SCRIPT_PATH/helloworld_script.sh
